# deploy/user-level/appsflyer-daily.service
[Unit]
Description=AppsFlyer -> analytics DB daily incremental load (BAF-2) [user-level stopgap, no root]
Documentation=https://yesimapp.atlassian.net/browse/BAF-2
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
# See docs/RUNBOOK.md for why RemainAfterExit stays off on a timer-triggered oneshot.

WorkingDirectory=%h/GitHubRepos/appsflyer-to-analytics-db
EnvironmentFile=%h/appsflyer-secrets/appsflyer.env
Environment=PYTHONUNBUFFERED=1 PYTHONDONTWRITEBYTECODE=1

ExecStart=%h/GitHubRepos/appsflyer-to-analytics-db/.venv/bin/appsflyer-pipeline daily

TimeoutStartSec=1800

StandardOutput=journal
StandardError=journal
SyslogIdentifier=appsflyer-daily

# Reduced hardening vs. ../appsflyer-daily.service: this is a systemd --user
# (per-user manager) unit, run without root -- see docs/RUNBOOK.md's
# "No-root stopgap deployment" section for why and when to retire this.
# ProtectHome/ProtectSystem=strict are intentionally omitted here: they would
# hide /home entirely, which is where the repo, venv, and secrets file live
# for a per-user deployment.
# ProtectClock/ProtectKernelModules/ProtectKernelLogs are deliberately ABSENT
# here (the root-based unit keeps them): each implies dropping capabilities
# from the bounding set, which an unprivileged user manager can only do inside
# a user namespace — and Ubuntu 24.04's default
# kernel.apparmor_restrict_unprivileged_userns=1 forbids that, failing the
# unit with 218/CAPABILITIES before the pipeline even starts (issue #19,
# observed on the first scheduled fire). They'd be security no-ops here
# anyway: an unprivileged process never holds CAP_SYS_TIME/SYS_MODULE/SYSLOG.
NoNewPrivileges=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
ProtectHostname=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
UMask=0077

# No [Install] section here either -- the .timer pulls this in.
